# shared/consul/init/systemd/consul-online.target
[Unit]
Description=Consul Online
RefuseManualStart=true
