# shared/nomad/init/systemd/nomad-online.target
[Unit]
Description=Nomad Online
RefuseManualStart=true
