# shared/nomad/init/systemd/vault-token-ready.target
[Unit]
Description=Vault Introduction Complete
RefuseManualStart=true
